# containers/rootfs/etc/s6-overlay/s6-rc.d/init-data/up
#!/command/execlineb -P
with-contenv
run-parts /etc/s6-overlay/init-data/