# containers/rootfs/etc/s6-overlay/s6-rc.d/chrome/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck -d
s6-setuidgid zeroclaw
chromedriver --port=9515 --allowed-ips=127.0.0.1