# containers/rootfs/etc/s6-overlay/s6-rc.d/zeroclaw/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck -d
s6-setuidgid zeroclaw
zeroclaw daemon